# test_unit.service
[Unit]
Description=Test unit script

[Service]
ExecStart=/home/taufiq/test_script

[Install]
WantedBy=multi-user.target
